# deploy/rpi/Services/CameraThumbnails/CheckThumbnails.timer
[Unit]
Description=OSH Create Thumbnail Timer

[Timer]
OnUnitActiveSec=12h
OnBootSec=12h

[Install]
WantedBy=timers.target